# roles/hashivault_deploy/templates/vault.service.j2
# Generated by Ansible

[Unit]
Description=HashiCorp Vault - A tool for managing secrets
Documentation=https://developer.hashicorp.com/vault/docs
Requires=network-online.target
After=network-online.target
ConditionFileNotEmpty={{ hashivault_deploy_config_dir }}/vault.hcl
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=notify
User={{ hashivault_deploy_user }}
Group={{ hashivault_deploy_group }}
EnvironmentFile={{ hashivault_deploy_config_dir }}/vault.env
ExecStart=/usr/bin/vault server -config={{ hashivault_deploy_config_dir }}/vault.hcl
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGINT
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536
LimitMEMLOCK=infinity
LimitCORE=0

{% if hashivault_deploy_systemd_extras %}
# Extra systemd security settings
{%   for key, value in hashivault_deploy_systemd_content.items() %}
{{ key }}={{ value }}
{%   endfor %}
{% endif %}

[Install]
WantedBy=multi-user.target