# roles/hashivault_deploy/templates/vault.service.j2
# Generated by Ansible

[Unit]
Description=HashiCorp Vault - A tool for managing secrets
Documentation=https://www.vaultproject.io/docs/
Requires=network-online.target
ConditionFileNotEmpty={{ hashivault_deploy_config_dir }}/vault.hcl
StartLimitIntervalSec=60
StartLimitBurst=3


[Service]
Type=notify
User={{ hashivault_deploy_user }}
Group={{ hashivault_deploy_group }}
EnvironmentFile={{ hashivault_deploy_config_dir }}/vault.env
ExecStart=/usr/bin/vault server -config={{ hashivault_deploy_config_dir }}/vault.hcl
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGINT
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536
LimitMEMLOCK=infinity
LimitCORE=0

# Security
ProtectSystem=strict
ReadWritePaths=/opt/vault
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectProc=invisible
SecureBits=keep-caps
NoNewPrivileges=true
AmbientCapabilities={{ hashivault_deploy_capabilities
  | map('regex_replace', '(cap_[a-z_]+)=\\+ep', '\\1')
  | map('upper')
  | join(' ') }}

CapabilityBoundingSet={{ hashivault_deploy_capabilities
  | map('regex_replace', '(cap_[a-z_]+)=\\+ep', '\\1')
  | map('upper')
  | join(' ') }}

[Install]
WantedBy=multi-user.target